# udev/rules/96-rfkill.rules
SUBSYSTEM=="rfkill", RUN+="/usr/bin/su veids -c \"polybar-msg -p $(pgrep -f polybar\ top) hook airplane 1\""
