# udev/rules/96-rfkill.rules
SUBSYSTEM=="rfkill", RUN+="/usr/bin/su veids -c '/home/veids/.dotfiles/polybar/scripts/airplane.sh --trigger'"
